# templates/supervisor-nginx.conf
[program:nginx]
command=/usr/sbin/nginx -c /etc/nginx/nginx.conf

priority=999
autostart=true
autorestart=true

stdout_logfile=/var/log/supervisor/nginx-stdout.log
stderr_logfile=/var/log/supervisor/nginx-stderr.log
